# systemd/normal-checker.timer
[Unit]
Description=Normal Flight Checker Timer
Requires=normal-checker.service

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30sec

[Install]
WantedBy=timers.target
